# meshcore-decoder.service
[Unit]
Description=MeshCore packet decoder MQTT bridge
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/meshcore-tools
EnvironmentFile=/opt/meshcore-tools/decoder.env
ExecStart=/usr/bin/python3 /opt/meshcore-tools/decode.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
