# decoder/meshcore-decoder.service
[Unit]
Description=MeshCore packet decoder MQTT bridge
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=andy
WorkingDirectory=/opt/meshcore-tools/decoder
EnvironmentFile=/opt/meshcore-tools/decoder/decoder.env
ExecStart=/opt/meshcore-tools/venv/bin/python /opt/meshcore-tools/decoder/decode.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
